# src/ansible/roles/setup/files/mount/rp02/srv.automount
[Unit]
Description=Automount /srv

[Automount]
Where=/srv
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
